# docker/datomic/datomic.conf
[program:datomic]
command=/datomic/bin/transactor /datomic/transactor.properties
stderr_logfile = /var/log/supervisord/datomic-stderr.log
stdout_logfile = /var/log/supervisord/datomic-stdout.log
