# docker/orcpub/orcpub.conf
[program:orcpub]
command=java -jar /orcpub.jar
stderr_logfile = /var/log/supervisord/orcpub-stderr.log
stdout_logfile = /var/log/supervisord/orcpub-stdout.log
